# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

[program:api]
command=python -m nginx_manager
directory=/app
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
stdout_logfile=/var/log/api.log
stderr_logfile=/var/log/api.log
environment=DATABASE_URL=sqlite:////app/data/nginx_manager.db,SECRET_KEY=%(ENV_SECRET_KEY)s,LOG_LEVEL=%(ENV_LOG_LEVEL)s

[program:mcp]
command=/app/.venv/bin/python -m nginx_manager.mcp
directory=/app
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
stdout_logfile=/var/log/mcp.log
stderr_logfile=/var/log/mcp.log
environment=DATABASE_URL=sqlite:////app/data/nginx_manager.db,SECRET_KEY=%(ENV_SECRET_KEY)s,LOG_LEVEL=%(ENV_LOG_LEVEL)s
